# b2_backups/duplicity/systemd/duplicity-backup.timer
[Unit]
Description=Run duplicity once a week

[Timer]
OnCalendar=Tue *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
