# b2_backups/duplicity/systemd/duplicity-backup.service
[Unit]
Description=Run duplicity backup

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'duplicity-backup.sh backup'